# scripts/systemd/zaoscribe.service
[Unit]
Description=ZAOscribe - Discord/Telegram meeting capture bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=zaoscribe
Group=zaoscribe
WorkingDirectory=/opt/zaoscribe
EnvironmentFile=/opt/zaoscribe/.env
ExecStart=/usr/bin/node --import tsx /opt/zaoscribe/src/index.ts
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Hardening - app dir read-only, recordings live in /var/lib so disk usage is
# observable + backup-policy-distinct from the code dir. Make sure
# RECORDINGS_DIR in .env matches the ReadWritePaths below.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/var/lib/zaoscribe

[Install]
WantedBy=multi-user.target
